# wsdl/autostart/oc4j.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<!--
        Copyright 2006,2009 Yuri Voinov. All rights reserved.
        Service manifest for OC4J services.      Version 1.7
-->

<service_bundle type='manifest' name='oracle-services'>

<service
        name='application/oc4j'
        type='service'
        version='1'>

        <create_default_instance enabled='false' />

        <single_instance />

        <dependency name="multi-user" 
                grouping="require_all" 
              restart_on="error" 
                    type="service">
           <service_fmri value="svc:/milestone/multi-user:default"/>
        </dependency>

        <dependency name="oracle" 
                grouping="require_all" 
              restart_on="error" 
                    type="service">
           <service_fmri value="svc:/application/oracle:default"/>
        </dependency>

        <exec_method
                type='method'
                name='start'
                exec='/lib/svc/method/init.oc4j %m'
                timeout_seconds='0' />

        <exec_method
                type='method'
                name='stop'
                exec='/lib/svc/method/init.oc4j %m'
                timeout_seconds='0' />

        <exec_method
                type='method'
                name='restart'
                exec='/lib/svc/method/init.oc4j %m'
                timeout_seconds='0' />

        <property_group name='options' type='application'>
                <propval name='tcp_listen' type='boolean'
                        value='true' />
        </property_group>

        <stability value='External' />

        <template>
                <common_name>
                        <loctext xml:lang='C'> 
                         Oracle Standalone OC4J Services 
                        </loctext>
                </common_name>
        </template>
</service>

</service_bundle>